# .cursor-server/data/User/History/-2d327a67/o72s.service
[Unit]
Description=Zoho-Odoo Continuous Product Sync Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/odoo/migration
ExecStart=/usr/bin/python3 /opt/odoo/migration/continuous_sync_server.py
Restart=always
RestartSec=30
StandardOutput=append:/opt/odoo/migration/logs/sync_service.log
StandardError=append:/opt/odoo/migration/logs/sync_service_error.log

# Environment variables
Environment=PYTHONPATH=/opt/odoo/migration
Environment=SYNC_INTERVAL=300

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/odoo/migration

[Install]
WantedBy=multi-user.target 